# get_config_files.list

### OS
/etc/fstab
/etc/group
/etc/passwd
/etc/sysconfig/network-scripts/ifcfg-eth0
/etc/sysconfig/network-scripts/ifcfg-eth1
/etc/sysconfig/network
/etc/modprobe.d/disable-ipv6.conf
/etc/hosts
/etc/resolv.conf
/etc/ssh/sshd_config
/etc/inittab
/etc/ntp.conf
/etc/sysconfig/selinux
/etc/crontab

### Apache 
/etc/httpd/conf/httpd.conf
/etc/httpd/conf.d/httpd-vhost-epit.conf

### php 
/etc/php.ini

### sendmail 
/etc/mail/sendmail.mc
/etc/mail/dkim-hosts
/etc/mail/dkim-keys
/etc/mail/dkim-milter/dkim-filter.conf
/etc/mail/dk-milter/devtools/Site/site.config.m4
/etc/rc.d/init.d/dk-filter

### Tomcat 
/usr/local/apache-tomcat-7.0.52/conf/server.xml

### MySQL
/etc/my.cnf

### dovecot
/etc/dovecot/dovecot.conf
/etc/dovecot/conf.d/10-mail.conf
/etc/dovecot/conf.d/10-auth.conf
/etc/mail/local-host-names
/etc/mail/virtusertable
/etc/aliases



